# ir_keytable.service
[Unit]
Description=Load ir-ketable definition
Requires=multi-user.target
Wants=systemd-modules-load.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=ir-keytable -c -w "/etc/rc_keymaps/samsung_necx.toml"
ExecReload=ir-keytable -c -w "/etc/rc_keymaps/samsung_necx.toml"
ExecStop=ir-keytable -c

[Install]
WantedBy=multi-user.target
