# ir_keytable.service
[Unit]
Description=ir-keytable

[Service]
ExecStart=ir-keytable -c -w /etc/rc_keymaps/samsung_necx.toml

[Install]
WantedBy=multi-user.target
